# pet_monitor/pet-monitor.service
[Unit]
Description=Pet Monitor Audio Service
After=network.target

[Service]
Type=simple
User=adam
WorkingDirectory=/home/adam/pet-monitor
Environment=PYTHONUNBUFFERED=1
Environment=ENV=production
ExecStart=/home/adam/pet-monitor/venv/bin/python audio_monitor.py
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
